# my_files/files/etc/uci-defaults/9999-my-custom-defaults.sh
#!/bin/sh

sed -i '42d;43d;44d;45d;53d;54d' /etc/board.json

exit 0


